# alpha-quant/supervisor.conf
; supervisor.conf --- 进程守护配置
; 管理所有Alpha Quant服务

[supervisord]
logfile=/root/.openclaw/workspace/alpha-quant/logs/supervisord.log
logfile_maxbytes=10MB
logfile_backups=3
loglevel=info
pidfile=/tmp/supervisord.pid
nodaemon=false

; 信号服务器
[program:signal_server]
command=/root/.openclaw/workspace/alpha-quant/quant_env/bin/python /root/.openclaw/workspace/alpha-quant/signal_server.py
directory=/root/.openclaw/workspace/alpha-quant
autostart=true
autorestart=true
startsecs=5
startretries=3
stderr_logfile=/root/.openclaw/workspace/alpha-quant/logs/signal_server_err.log
stdout_logfile=/root/.openclaw/workspace/alpha-quant/logs/signal_server_out.log
environment=SIGNAL_PORT="8765",FEISHU_WEBHOOK_URL="https://open.feishu.cn/open-apis/bot/v2/hook/ba2ce0cf-7d1a-43aa-ad93-d8747b8a24b0"

; 数据网关（使用Gunicorn）
[program:data_gateway]
command=/root/.openclaw/workspace/alpha-quant/quant_env/bin/gunicorn -w 2 -b 0.0.0.0:8766 --timeout 60 data_gateway:app
directory=/root/.openclaw/workspace/alpha-quant
autostart=true
autorestart=true
startsecs=5
startretries=3
stderr_logfile=/root/.openclaw/workspace/alpha-quant/logs/gateway_err.log
stdout_logfile=/root/.openclaw/workspace/alpha-quant/logs/gateway_out.log

; 统一调度器
[program:scheduler]
command=/root/.openclaw/workspace/alpha-quant/quant_env/bin/python /root/.openclaw/workspace/alpha-quant/scheduler.py
directory=/root/.openclaw/workspace/alpha-quant
autostart=true
autorestart=true
startsecs=10
startretries=3
stderr_logfile=/root/.openclaw/workspace/alpha-quant/logs/scheduler_err.log
stdout_logfile=/root/.openclaw/workspace/alpha-quant/logs/scheduler_out.log

; 组管理
[group:alpha_quant]
programs=signal_server,data_gateway,scheduler
